# maya-public/scripts/systemd/maya-ingest-poll.service
[Unit]
Description=Maya ingest subscription poll (YouTube Atom feeds → homepage notifications)
After=network-online.target postgresql.service
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=%h/Workspace-public
Environment=PGHOST=localhost
Environment=PGPORT=5433
Environment=PGUSER=maya
Environment=PGPASSWORD=maya
Environment=PGDATABASE=maya_public
Environment=DATABASE_URL=postgresql+asyncpg://maya:maya@localhost:5433/maya_public
ExecStart=/usr/bin/make ingest-poll

[Install]
WantedBy=default.target
