# scripts/systemd/maya-gateway.service
[Unit]
Description=Maya gateway (hyprstart homepage + API on :8090)
After=network-online.target postgresql.service
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/Workspace-public
Environment=PGHOST=localhost
Environment=PGPORT=5433
Environment=PGUSER=maya
Environment=PGPASSWORD=maya
Environment=PGDATABASE=maya_public
Environment=DATABASE_URL=postgresql+asyncpg://maya:maya@localhost:5433/maya_public
Environment=PORT=8090
Environment=ENV=production
ExecStart=/usr/bin/make gateway-dev
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
